# vagrant_provision/supervisor.default.conf
[program:site]
command=/home/vagrant/site/bin/python /vagrant/site/site.py
